# newapi-guardian.service.example
[Unit]
Description=NewAPI Guardian Telegram Bot
After=network.target mysql.service

[Service]
Type=simple
User=your_user
WorkingDirectory=/path/to/newapi-guardian
Environment="PATH=/path/to/newapi-guardian/.venv/bin"
ExecStart=/path/to/newapi-guardian/.venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
